# _config/ipython_notebook.conf
; supervisord configuration file for ipython notebook
[program:ipython_notebook]
command=ipython notebook --ip 0.0.0.0
autostart=true
autorestart=true
directory=/vagrant